# scripts/mealie-mcp-fraser.service
# Systemd user service for Mealie MCP Server (Fraser's instance)
# 
# Installation:
#   1. Copy to ~/.config/systemd/user/mealie-mcp-fraser.service
#   2. Run: systemctl --user daemon-reload
#   3. Run: systemctl --user enable --now mealie-mcp-fraser

[Unit]
Description=Mealie MCP Server (Fraser)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/fraser/Repos/mealie-mcp-server
Environment="PATH=/home/fraser/Repos/mealie-mcp-server/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="MCP_TRANSPORT=sse"
Environment="MCP_HOST=0.0.0.0"
Environment="MCP_PORT=8080"
Environment="MCP_USER=fraser"
Environment="MCP_BASE_URL=https://rainworth-server.tailbf31d9.ts.net/fraser"
EnvironmentFile=/home/fraser/Repos/mealie-mcp-server/.env
ExecStart=/home/fraser/Repos/mealie-mcp-server/.venv/bin/python -m mealie_mcp.server
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
